# supervisor_sample.conf
# [Supervisor Sample Configuration for VPN Panel]
# Lưu tại: /etc/supervisor/conf.d/vpn_panel.conf

[group:vpn]
programs=vpn_panel,vpn_worker,vpn_beat

# 1. Giao diện Web (Gunicorn)
[program:vpn_panel]
directory=/data/openvpn
command=/home/kien/miniconda3/bin/gunicorn vpn_project.wsgi:application --bind 127.0.0.1:8003 --workers 3 --timeout 120
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/vpn_panel.out.log
stderr_logfile=/var/log/vpn_panel.err.log
environment=LANG=en_US.UTF_8,LC_ALL=en_US.UTF_8
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
stopasgroup=true

# 2. Celery Worker (Xử lý tác vụ nền & Locking)
[program:vpn_worker]
command=/home/kien/miniconda3/bin/celery -A vpn_project worker -l info
directory=/data/openvpn
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/vpn_worker.out.log
stderr_logfile=/var/log/vpn_worker.err.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
stopasgroup=true

# 3. Celery Beat (Lập lịch tự động dọn dẹp & Thông báo)
[program:vpn_beat]
command=/home/kien/miniconda3/bin/celery -A vpn_project beat -l info --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/data/openvpn
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/vpn_beat.out.log
stderr_logfile=/var/log/vpn_beat.err.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
stopasgroup=true

# --- Ghi chú Cài đặt ----------------------------------------------------
# 1. Thay đổi 'user=root' thành user thực tế nếu không muốn chạy quyền cao.
# 2. Kiểm tra đường dẫn 'command' và 'directory' đúng với server của bạn.
# 3. Sau khi cập nhật file này, chạy:
#    sudo supervisorctl reread
#    sudo supervisorctl update
#    sudo supervisorctl restart vpn:*
# -----------------------------------------------------------------------
